# scripts/bellforge-client.service
[Unit]
Description=BellForge Signage Client (Chromium Kiosk)
Documentation=https://github.com/YOUR_ORG/BellForge
After=graphical.target network-online.target bellforge-backend.service
Wants=network-online.target bellforge-backend.service

[Service]
Type=simple
User=bellforge
Group=bellforge

Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/bellforge/.Xauthority
Environment=BELLFORGE_KIOSK_URL=http://127.0.0.1:8000/status
EnvironmentFile=-/opt/bellforge/config/client.env

# Open the status dashboard in full-screen kiosk mode.
ExecStart=/bin/sh -lc 'BROWSER="${BELLFORGE_BROWSER_BIN:-$(command -v chromium-browser || command -v chromium)}"; exec "$BROWSER" --kiosk --noerrdialogs --disable-infobars --disable-session-crashed-bubble --disable-restore-session-state --disable-background-networking --disable-translate --app="${BELLFORGE_KIOSK_URL}"'

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bellforge-client

[Install]
WantedBy=graphical.target
