# scripts/bellforge-client.service
[Unit]
Description=BellForge Signage Client (Chromium Kiosk)
Documentation=https://github.com/ronaldarroyowatson/BellForge
After=graphical.target network-online.target bellforge-backend.service
Wants=network-online.target bellforge-backend.service

[Service]
Type=simple
User=bellforge
Group=bellforge

Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/bellforge/.Xauthority
Environment=BELLFORGE_KIOSK_URL=http://127.0.0.1:8000/client/index.html
EnvironmentFile=-/opt/bellforge/config/client.env

# Open the status dashboard in full-screen kiosk mode.
ExecStart=/opt/bellforge/scripts/start_kiosk.sh

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bellforge-client

[Install]
WantedBy=graphical.target
